# deployment/systemd/vasi-engine-egress-boundary.service
[Unit]
Description=Verify the VASI private-engine egress boundary
Documentation=file:/opt/vasi-engine/current/docs/architecture/private-engine-egress.md
Requires=docker.service
After=docker.service network-online.target vasi-engine-database-egress-policy.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
WorkingDirectory=/opt/vasi-engine/current
ExecStart=/usr/bin/env node scripts/probe-engine-egress-boundary.mjs
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
LockPersonality=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictSUIDSGID=yes
SystemCallArchitectures=native
