# deployment/systemd/vasi-engine-database-egress-policy.timer
[Unit]
Description=Refresh the VASI engine exact database and private-ingress egress policy

[Timer]
OnBootSec=15s
OnUnitActiveSec=2min
AccuracySec=15s
RandomizedDelaySec=15s
Persistent=yes
Unit=vasi-engine-database-egress-policy.service

[Install]
WantedBy=timers.target
