# deployment/systemd/vasi-engine-alert-readiness.timer
[Unit]
Description=Continuously expose pending VASI private-engine operational alerts

[Timer]
OnBootSec=30s
OnActiveSec=30s
OnUnitInactiveSec=1min
AccuracySec=15s
RandomizedDelaySec=15s
Persistent=yes
Unit=vasi-engine-alert-readiness.service

[Install]
WantedBy=timers.target
