# deployment/systemd/vasi-engine-operational-readiness.service
[Unit]
Description=Verify privacy-safe VASI private-engine operational readiness
Documentation=file:/opt/vasi-engine/current/docs/architecture/operational-readiness.md
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
WorkingDirectory=/opt/vasi-engine/current
ExecStart=/usr/bin/docker compose -f compose.engine.yaml --profile tools run --rm -T --no-deps maintenance scripts/probe-operational-readiness.mjs
TimeoutStartSec=5min
Nice=10
IOSchedulingClass=idle
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictSUIDSGID=yes
SystemCallArchitectures=native
