# deployment/systemd/vasi-engine-backup-create.timer
[Unit]
Description=Schedule verified matched VASI private-engine backups

[Timer]
OnBootSec=30min
OnActiveSec=30min
OnUnitInactiveSec=24h
AccuracySec=15min
RandomizedDelaySec=30min
Persistent=yes
Unit=vasi-engine-backup-create.service

[Install]
WantedBy=timers.target
