# deployment/systemd/vasi-engine-database-egress-policy.service
[Unit]
Description=VASI engine exact database and private-ingress egress policy
Documentation=file:/opt/vasi-engine/current/docs/architecture/private-engine-egress.md
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
OnFailure=vasi-engine-alert-record@%n.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
WorkingDirectory=/opt/vasi-engine/current
ExecStart=/bin/sh scripts/apply-database-egress-policy.sh apply
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/run
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictSUIDSGID=yes
SystemCallArchitectures=native
