# deployment/systemd/vasi-edge-alert-readiness.service
[Unit]
Description=Verify the durable VASI public-edge operational-alert handoff
Documentation=file:/opt/vasi-edge/current/docs/architecture/durable-operational-alert-handoff.md

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
ExecStart=/bin/sh /usr/local/libexec/vasi/operational-alert-spool.sh status edge
TimeoutStartSec=30s
StateDirectory=vasi-edge/operations-alerts
StateDirectoryMode=0700
CapabilityBoundingSet=
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
InaccessiblePaths=/run
ReadWritePaths=/var/lib/vasi-edge/operations-alerts
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
