# deployment/systemd/vasi-engine-backup-create.service
[Unit]
Description=Create a verified matched VASI private-engine backup
Documentation=file:/opt/vasi-engine/current/docs/architecture/backup-continuity.md
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
OnFailure=vasi-engine-alert-record@%n.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
WorkingDirectory=/opt/vasi-engine/current
ExecStart=/usr/bin/docker compose -f compose.engine.yaml --profile tools run --rm -T --no-deps -v /var/lib/vasi-engine/backups/maintenance/scheduled:/backup maintenance scripts/backup-continuity.mjs create /backup
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=idle
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictSUIDSGID=yes
SystemCallArchitectures=native
